# tg_bot_updater.service
[Unit]
Description=Telegram Bot Auto Updater Service
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/TG_BOT_TU-UGMK
ExecStart=/bin/bash /root/TG_BOT_TU-UGMK/auto_updater.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
